# Chapter-18/photogallery.conf
[program:photogallery]
directory=/home/ubuntu/photogallery/
command=gunicorn --bind 0.0.0.0:5000 app
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true

